# crontab
* * * * * root python3 /app/youtube_dl_script.py >> /proc/1/fd/1 2>> /proc/1/fd/2
